# sample_codes_updated/7-2/crawler.timer
[Unit]
Description=Timer for python crawler

[Timer]
OnCalendar=6:13
Unit=crawler.service

[Install]
WantedBy=multi-user.target
