# sample_codes_updated/7-2/crawler.service
[Unit]
Description=Python crawler

[Service]
Type=oneshot
ExecStart=/home/ubuntu/crawler/run_crawler.sh
User=ubuntu
